# com.user.clipboardmanager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.clipboardmanager</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/ClipboardManager.app/Contents/MacOS/ClipboardManager</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clipboardmanager.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/clipboardmanager.out</string>
</dict>
</plist> 